# repo-watcher.service
[Unit]
Description=Repo Watcher Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/repo-watcher/multi_monitor.py
Restart=always
RestartSec=10
WorkingDirectory=/opt/repo-watcher
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target